# sysmon-tracker.service
[Unit]
Description=System Monitor Screen Time Tracker
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/sysmon-tracker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
